# dockerCustomImages/kaliDockerDesktop/supervisord.conf
[supervisord]
nodaemon=true  ; Ensures supervisord runs in the foreground

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman
autostart=true
autorestart=true
stderr_logfile=/var/log/xrdp-sesman.err.log
stdout_logfile=/var/log/xrdp-sesman.out.log
user=root  ; xrdp-sesman needs to run as root

[program:xrdp]
command=/usr/sbin/xrdp
autostart=true
autorestart=true
stderr_logfile=/var/log/xrdp.err.log
stdout_logfile=/var/log/xrdp.out.log
user=root  ; xrdp needs to run as root
